# ze-scripts.supervisor.conf
[program:ze-v1]
command=/usr/bin/node /app/zedelivery-clean/puppeteer-wrapper.js v1.js
directory=/app/zedelivery-clean
autostart=true
autorestart=true
startretries=999
startsecs=5
stopwaitsecs=10
stdout_logfile=/app/logs/ze-v1-out.log
stderr_logfile=/app/logs/ze-v1-error.log
environment=NODE_ENV="production",PUPPETEER_EXECUTABLE_PATH="/usr/bin/chromium"

[program:ze-v1-itens]
command=/usr/bin/node /app/zedelivery-clean/puppeteer-wrapper.js v1-itens.js
directory=/app/zedelivery-clean
autostart=true
autorestart=true
startretries=999
startsecs=5
stopwaitsecs=10
stdout_logfile=/app/logs/ze-v1-itens-out.log
stderr_logfile=/app/logs/ze-v1-itens-error.log
environment=NODE_ENV="production",PUPPETEER_EXECUTABLE_PATH="/usr/bin/chromium"

[program:ze-sync]
command=/usr/bin/node /app/bridge/sync-cron.js
directory=/app/bridge
autostart=true
autorestart=true
startretries=999
startsecs=5
stopwaitsecs=10
stdout_logfile=/app/logs/ze-sync-out.log
stderr_logfile=/app/logs/ze-sync-error.log
environment=NODE_ENV="production"
