# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ============================================
; Backend FastAPI - Serve API + Frontend
; ÚNICO PROCESSO WEB - Railway precisa disso
; ============================================
[program:backend]
command=sh -c "python3 -m uvicorn server:app --host 0.0.0.0 --port ${PORT:-8080}"
directory=/app/backend
autostart=true
autorestart=true
priority=1
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

; ============================================
; Sync para Lovable Cloud (worker em background)
; ============================================
[program:ze-sync]
command=node /app/bridge/sync-cron.js
directory=/app/bridge
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=10
stdout_logfile=/app/logs/ze-sync-out.log
stderr_logfile=/app/logs/ze-sync-err.log
environment=NODE_ENV="production"

; ============================================
; SCRAPERS DESABILITADOS TEMPORARIAMENTE
; Puppeteer/Chromium precisa de configuração especial no Railway
; Eles continuam rodando no Emergent/VPS atual
; ============================================
; [program:ze-v1]
; command=node /app/zedelivery-clean/v1.js
; directory=/app/zedelivery-clean
; autostart=false
; 
; [program:ze-v1-itens]
; command=node /app/zedelivery-clean/v1-itens.js
; directory=/app/zedelivery-clean
; autostart=false
