# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; Healthcheck PRIMEIRO - responde na porta do Railway
[program:healthcheck]
command=node /app/healthcheck-server.js
directory=/app
autostart=true
autorestart=true
priority=1
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Scraper principal - lista de pedidos
[program:ze-v1]
command=node /app/zedelivery-clean/v1.js
directory=/app/zedelivery-clean
autostart=true
autorestart=true
startretries=10
startsecs=10
priority=10
stdout_logfile=/app/logs/ze-v1-out.log
stderr_logfile=/app/logs/ze-v1-err.log
environment=NODE_ENV="production",PUPPETEER_SKIP_CHROMIUM_DOWNLOAD="true",PUPPETEER_EXECUTABLE_PATH="/usr/bin/chromium"

; Scraper de itens - detalhes dos pedidos
[program:ze-v1-itens]
command=node /app/zedelivery-clean/v1-itens.js
directory=/app/zedelivery-clean
autostart=true
autorestart=true
startretries=10
startsecs=10
priority=10
stdout_logfile=/app/logs/ze-v1-itens-out.log
stderr_logfile=/app/logs/ze-v1-itens-err.log
environment=NODE_ENV="production",PUPPETEER_SKIP_CHROMIUM_DOWNLOAD="true",PUPPETEER_EXECUTABLE_PATH="/usr/bin/chromium"

; Sync para Lovable Cloud
[program:ze-sync]
command=node /app/bridge/sync-cron.js
directory=/app/bridge
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=10
stdout_logfile=/app/logs/ze-sync-out.log
stderr_logfile=/app/logs/ze-sync-err.log
environment=NODE_ENV="production"
