# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:ze-v1]
command=node /app/zedelivery-clean/v1.js
directory=/app/zedelivery-clean
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/app/logs/ze-v1-out.log
stderr_logfile=/app/logs/ze-v1-err.log
environment=NODE_ENV="production"

[program:ze-v1-itens]
command=node /app/zedelivery-clean/v1-itens.js
directory=/app/zedelivery-clean
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/app/logs/ze-v1-itens-out.log
stderr_logfile=/app/logs/ze-v1-itens-err.log
environment=NODE_ENV="production"

[program:ze-sync]
command=node /app/bridge/sync-cron.js
directory=/app/bridge
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/app/logs/ze-sync-out.log
stderr_logfile=/app/logs/ze-sync-err.log
environment=NODE_ENV="production"

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -F
autostart=true
autorestart=true
stdout_logfile=/app/logs/php-fpm-out.log
stderr_logfile=/app/logs/php-fpm-err.log

[program:healthcheck]
command=node /app/healthcheck-server.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/healthcheck-out.log
stderr_logfile=/app/logs/healthcheck-err.log
environment=PORT="%(ENV_PORT)s"
