# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ============================================
; Backend FastAPI - Serve API + Frontend
; ============================================
[program:backend]
command=sh -c "python3 -m uvicorn server:app --host 0.0.0.0 --port ${PORT:-8080}"
directory=/app/backend
autostart=true
autorestart=true
priority=1
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

; ============================================
; Scraper principal
; ============================================
[program:ze-v1]
command=node v1.js
directory=/app/zedelivery-clean
autostart=true
autorestart=true
startretries=5
startsecs=15
priority=10
stdout_logfile=/app/logs/ze-v1-out.log
stderr_logfile=/app/logs/ze-v1-err.log
environment=NODE_ENV="production",PUPPETEER_EXECUTABLE_PATH="/usr/bin/chromium"

; ============================================
; Scraper de itens
; ============================================
[program:ze-v1-itens]
command=node v1-itens.js
directory=/app/zedelivery-clean
autostart=true
autorestart=true
startretries=5
startsecs=15
priority=10
stdout_logfile=/app/logs/ze-v1-itens-out.log
stderr_logfile=/app/logs/ze-v1-itens-err.log
environment=NODE_ENV="production",PUPPETEER_EXECUTABLE_PATH="/usr/bin/chromium"

; ============================================
; Sync para Lovable Cloud
; ============================================
[program:ze-sync]
command=node sync-cron.js
directory=/app/bridge
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=10
stdout_logfile=/app/logs/ze-sync-out.log
stderr_logfile=/app/logs/ze-sync-err.log
environment=NODE_ENV="production"
